# package/dist/lib/systemd/system/hyperd.service
[Unit]
Description=hyperd
Documentation=http://docs.hyper.sh
After=network.target
Requires=

[Service]
ExecStart=/usr/bin/hyperd
MountFlags=slave
LimitNOFILE=1048576
LimitNPROC=1048576
LimitCORE=infinity

[Install]
WantedBy=multi-user.target
